# .config/autostart/spotifyd.desktop
[Desktop Entry]
Name=Spotifyd
Exec=spotifyd --no-daemon --autoplay -b pulseaudio
Type=Application
