# .config/autostart/spotifyd.desktop
[Desktop Entry]
Name=Spotifyd
Exec=spotifyd --no-daemon --autoplay --backend pulseaudio
Type=Application
